# check_sync.crontab
# CNCF devstats.cncf.io K8s node-0: check devstats landscape sync
0 3 * * * CHECK_SYNC_DIR=/root/go/src/github.com/cncf/devstats-landscape-sync GOPATH=/root/go GOROOT=/usr/local/go PATH=$PATH:/usr/sbin:/root/go/bin/:/usr/local/go/bin /root/go/src/github.com/cncf/devstats-landscape-sync/check_sync.sh 1>>/tmp/check_sync.log 2>>/tmp/check_sync.err
